# telegram-crypto-alerts/xe-price-alert-bot.service
[Unit]
Description=Xe Price Alert Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/xe-price-alert-bot
EnvironmentFile=/opt/xe-price-alert-bot/.env
ExecStart=/opt/xe-price-alert-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
